# packages/wooting-mouse/wooting_mouse-0.4.0-py3-none-any.whl/data/wooting-mouse@.service
[Unit]
Description=Wooting Mouse daemon

[Service]
Type=simple
ExecStart=/usr/bin/wooting-mouse /%I
StandardOutput=journald
StandardError=journald
Environment="PYTHONUNBUFFERED=1"
